# talevision.service
[Unit]
Description=TaleVision e-ink display appliance
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/talevision
ExecStart=/home/pi/talevision/venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=talevision
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
